# locking/bash/supervisord.conf
[supervisord]
nodaemon=true

[program:bash]
command=/usr/bin/consul-lock.bash

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
